# data/system/eos-kolibri-system-helper.service.in
[Unit]
Description=Kolibri Server
ConditionPathExists=/var/lib/flatpak/app/org.learningequality.Kolibri

[Service]
Type=simple
ExecStart=@bindir@/flatpak run --no-desktop --filesystem=@KOLIBRI_DATA_DIR@ --command=kolibri org.learningequality.Kolibri start --foreground
Environment=KOLIBRI_HOME=@KOLIBRI_DATA_DIR@
Environment=KOLIBRI_HTTP_PORT=@KOLIBRI_PORT@
Environment=KOLIBRI_IDLE_TIMEOUT_MINS=10
User=@KOLIBRI_USER@
PrivateTmp=yes
Restart=on-failure
