# data/system/eos-kolibri-system-helper.socket.in
[Unit]
Description=Kolibri Server

[Socket]
ListenStream=__KOLIBRI_PORT__
Accept=no

[Install]
WantedBy=sockets.target
