# data/systemd/dbus-org.learningequality.Kolibri.Daemon.service.in
[Unit]
Description=Kolibri Daemon Service
ConditionPathExists=/var/lib/flatpak/app/@KOLIBRI_FLATPAK_ID@

[Service]
Type=dbus
BusName=@KOLIBRI_DAEMON_SERVICE@
ExecStart=@libexecdir@/@KOLIBRI_CMD_PREFIX@-daemon
User=@KOLIBRI_USER@
PrivateTmp=yes
