# data/systemd/dbus-org.learningequality.Kolibri.Daemon.service.in
[Unit]
Description=Kolibri Daemon Service
ConditionPathExists=/var/lib/flatpak/app/@KOLIBRI_FLATPAK_ID@

[Service]
Type=dbus
BusName=@KOLIBRI_DAEMON_SERVICE@
ExecStart=@libexecdir@/eos-kolibri-daemon
DynamicUser=yes
User=@KOLIBRI_USER@
RuntimeDirectory=kolibri
StateDirectory=kolibri
Environment=HOME=%t/kolibri
